# src/static/static/home/yi-hack/etc/crontabs/root
0 * * * * /home/yi-hack/script/check_update.sh
